# deploy/aliyun/minasan.service
[Unit]
Description=Minasan PWA Node server
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/minasan
Environment=NODE_ENV=production
Environment=HOST=127.0.0.1
Environment=PORT=8788
Environment=PUBLIC_ORIGIN=https://example.com
Environment=MINASAN_SQLITE_PATH=/var/lib/minasan/minasan.sqlite
ExecStart=/usr/bin/npm run server:start
Restart=always
RestartSec=5
User=minasan
Group=minasan

[Install]
WantedBy=multi-user.target
